# luci-app-gpon/root/etc/uci-defaults/80_gpon
#!/bin/sh

touch /etc/config/gpon
uci set gpon.first=first
uci set gpon.second=second
uci set gpon.third=third
uci commit

return 0
